# zenoh/.deb/postinst
#!/bin/sh
# postinst script for Eclipse Zenoh
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
        id -u zenoh  >/dev/null 2>&1 ||  sudo useradd -r -s /bin/false zenoh
        mkdir -p /etc/zenoh
        chown zenoh:zenoh -R /etc/zenoh
        chmod 0644 /etc/zenoh -R
        systemctl daemon-reload
        systemctl disable zenoh
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0