# zenohd/.deb/postinst
#!/bin/sh
# postinst script for Eclipse Zenoh
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
        if ! command -v systemctl > /dev/null 2>&1
        then
            echo "WARNING: 'systemctl' not found - cannot install zenohd as a service."
            exit 0
        fi
        id -u zenohd > /dev/null 2>&1 || useradd -r -s /bin/false zenohd
        mkdir -p /var/zenohd
        chown zenohd:zenohd /var/zenohd
        systemctl daemon-reload
        systemctl disable zenohd
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0